# src/cli/attache.service.template
[Unit]
Description=Attache AI Assistant
After=network.target

[Service]
Type=simple
ExecStart={{ATTACHE_PATH}} start
Restart=always
RestartSec=5
Environment=PATH={{PATH}}
Environment=HOME={{HOME}}

[Install]
WantedBy=default.target
